# systemd/kidsnote.timer
[Unit]
Description=Run Kidsnote download script once per day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
